# external_ip_changed.path
# /usr/lib/systemd/system/external_ip_changed.path
[Unit]
Description=Trigger service if the IP changed
Requires=check_external_ip.timer
After=check_external_ip.timer

[Path]
PathModified=/tmp/external_ip

[Install]
WantedBy=multi-user.target
